# roles/airsonos/files/etc/systemd/system/airsonos.service
[Unit]
Description=AirSonos server

[Service]
Type=simple
User=root
ExecStart=/usr/bin/airsonos

[Install]
WantedBy=multi-user.target

